# systemd/system/iperf-plotter.timer
[Unit]
Description=Hourly measurement of network bandwidth
Documentation=

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h
